# com.vibevoice.tts-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.vibevoice.tts-server</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/hareesh/dev/ai/VibeVoice/start-tts-server.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/Users/hareesh/dev/ai/VibeVoice/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/hareesh/dev/ai/VibeVoice</string>

    <key>RunAtLoad</key>
    <true />

    <key>KeepAlive</key>
    <true />

    <key>StandardOutPath</key>
    <string>/Users/hareesh/dev/ai/VibeVoice/tts-server-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/hareesh/dev/ai/VibeVoice/tts-server-launchd.err</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>
  </dict>
</plist>